# ansible/roles/python_exporter/templates/guardian-exporter.service.j2
[Unit]
Description=Guardian custom Python exporter
After=network-online.target
Wants=network-online.target

[Service]
User=guardian-exporter
Group=guardian-exporter
Environment=GUARDIAN_EXPORTER_PORT={{ python_exporter_port | default(8000) }}
Environment=GUARDIAN_EXPORTER_STATE_FILE=/opt/guardian-exporter/state/sample_metrics.json
ExecStart=/opt/guardian-exporter/venv/bin/python -m guardian_exporter.app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
